# meta-custom/recipes-connectivity/ifup-wlan0/files/ifup-wlan0.service
[Unit]
Description=Bring up wlan0 at boot
After=network.target

[Service]
Type=oneshot
ExecStart=/sbin/ifup wlan0
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target